# etc/rc.local
#!/bin/bash
nohup python /home/pi/mail.py &
nohup python3 /home/pi/updateV2RAY.py &
